# services/mail/run
#!/bin/sh

while true
do
    mbsync -ac ~/.config/isync/mbsyncrc
    sleep 30
done
